# deploy/systemd/micro-kiki.service
[Unit]
Description=micro-kiki vLLM Serving
After=network.target

[Service]
Type=simple
User=kxkm
WorkingDirectory=/home/kxkm/micro-kiki
ExecStart=/home/kxkm/.local/bin/uv run python -m src.serving.vllm_server
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1
Environment=VLLM_ALLOW_RUNTIME_LORA_UPDATING=True

[Install]
WantedBy=multi-user.target
