# deploy/launchd/cc.saillant.micro-kiki.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>cc.saillant.micro-kiki</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/clems/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>src.serving.mlx_server</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/clems/micro-kiki</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/micro-kiki.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/micro-kiki.stderr.log</string>
</dict>
</plist>
